# 11_TP_Deployment_Ansible/13_TP_Deployment_Event_Driven_AWS_IoT_Core/Ansible/files/mqtt_listener.service
[Unit]
Description=MQTT Listener AWS
After=network.target

[Service]
User=vagrant
ExecStart=/home/vagrant/mqtt_listener.sh
Restart=always

[Install]
WantedBy=multi-user.target